# deploy/systemd/swing-scanner.service
[Unit]
Description=Swing Scanner one-shot run
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/swing_scanner
EnvironmentFile=/etc/swing-scanner.env
ExecStart=/opt/swing_scanner/.venv/bin/python -m swing_scanner.app --watchlist /opt/swing_scanner/watchlist.txt --run-once
User=swing
Group=swing
Nice=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
